# worker/agent-systemd-user.service
[Unit]
Description=systemd --user for agent (enables user-level services)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Create runtime dir for agent (UID 1001) so systemd --user can run
ExecStartPre=/bin/bash -c 'mkdir -p /run/user/1001 && chown agent:agent /run/user/1001'
# Pass container env into the user's systemd so openclaw.user.service can use them
PassEnvironment=HOME TERM PATH PNPM_HOME OPENCLAW_GATEWAY_BIND NEARAI_API_KEY OPENCLAW_GATEWAY_TOKEN OPENCLAW_FORCE_CONFIG_REGEN NODE_ENV
Environment=XDG_RUNTIME_DIR=/run/user/1001
Environment=OPENCLAW_GATEWAY_BIND=lan
ExecStart=/usr/sbin/runuser -p -u agent -- /lib/systemd/systemd --user
Restart=always
RestartSec=3
KillMode=mixed

[Install]
WantedBy=multi-user.target
