# worker/openclaw-gateway.service
[Unit]
Description=OpenClaw NEAR AI Gateway
After=network.target

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/home/agent
EnvironmentFile=-/home/agent/.openclaw/gateway.env
ExecStart=/home/agent/.npm-global/bin/openclaw gateway run --bind ${OPENCLAW_GATEWAY_BIND} --port 18789 --verbose
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
