# worker/openclaw.service
[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent
Group=agent
# Pass through env vars from container (set by docker/compose)
PassEnvironment=HOME TERM PATH PNPM_HOME OPENCLAW_GATEWAY_BIND NEARAI_API_KEY OPENCLAW_GATEWAY_TOKEN OPENCLAW_FORCE_CONFIG_REGEN NODE_ENV
Environment=OPENCLAW_GATEWAY_BIND=lan
ExecStart=/bin/bash -c 'exec /usr/sbin/runuser -p -u agent -- openclaw gateway run --bind "${OPENCLAW_GATEWAY_BIND}" --port 18789 --verbose'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target
