# worker/openclaw.user.service
[Unit]
Description=OpenClaw Gateway (user service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as agent; env inherited from systemd --user (set by agent-systemd-user.service)
Environment=OPENCLAW_GATEWAY_BIND=lan
ExecStart=/bin/bash -c 'exec openclaw gateway run --bind "${OPENCLAW_GATEWAY_BIND}" --port 18789 --verbose'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=default.target
